# lxd-panel/systemd/lxd-panel.service
[Unit]
Description=LXD Management Panel Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lxd-panel/backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/lxd-panel"
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
